# .claude/launchd/com.moonlight.buildsite.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moonlight.buildsite</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/zhenliu/Moonlight/.claude/scripts/build_site.py</string>
    </array>

    <!-- 每 30 分钟跑一次（1800 秒）-->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- 加载 plist 时立即跑一次 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 日志位置 -->
    <key>StandardOutPath</key>
    <string>/tmp/moonlight-buildsite.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/moonlight-buildsite.err</string>

    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <string>/Users/zhenliu/Moonlight</string>
</dict>
</plist>
